# com.hyperclaw.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hyperclaw</string>

    <key>ProgramArguments</key>
    <array>
        <string>~/.hyperclaw/.venv/bin/python3</string>
        <string>~/.hyperclaw/run_hyperclaw.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>~/.hyperclaw</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>~/.hyperclaw/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>~/.hyperclaw</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>~/.hyperclaw/logs/hyperclaw.log</string>

    <key>StandardErrorPath</key>
    <string>~/.hyperclaw/logs/hyperclaw.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
